# SupportFiles/LaunchAgents/com.github.execriez.labnotes.deploypayloads.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.github.execriez.labnotes.deploypayloads</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/LabNotes/bin/deploypayloads</string>
	</array>
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
	<key>StartInterval</key>
	<integer>179</integer>
</dict>
</plist>
